# contrib/init/hodlcoind.conf
description "HODLCoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env HODLCOIND_BIN="/usr/bin/hodlcoind"
env HODLCOIND_USER="hodlcoin"
env HODLCOIND_GROUP="hodlcoin"
env HODLCOIND_PIDDIR="/var/run/hodlcoind"
# upstart can't handle variables constructed with other variables
env HODLCOIND_PIDFILE="/var/run/hodlcoind/hodlcoind.pid"
env HODLCOIND_CONFIGFILE="/etc/hodlcoin/hodlcoin.conf"
env HODLCOIND_DATADIR="/var/lib/hodlcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # hodlcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$HODLCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run hodlcoind."
        echo "The setting must appear in $HODLCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"HODLCoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$HODLCOIND_PIDDIR"
    chmod 0755 "$HODLCOIND_PIDDIR"
    chown $HODLCOIND_USER:$HODLCOIND_GROUP "$HODLCOIND_PIDDIR"
    chown $HODLCOIND_USER:$HODLCOIND_GROUP "$HODLCOIND_CONFIGFILE"
    chmod 0660 "$HODLCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$HODLCOIND_PIDFILE" \
    --chuid $HODLCOIND_USER:$HODLCOIND_GROUP \
    --exec "$HODLCOIND_BIN" \
    -- \
    -pid="$HODLCOIND_PIDFILE" \
    -conf="$HODLCOIND_CONFIGFILE" \
    -datadir="$HODLCOIND_DATADIR" \
    -disablewallet \
    -daemon

